# robotics-lab/forgejo/robotics-lab-forgejo-backup.timer
[Unit]
Description=Nightly Nexus Robotics Forgejo backup

[Timer]
OnCalendar=*-*-* 03:15:00 America/Los_Angeles
Persistent=true
RandomizedDelaySec=15m
Unit=robotics-lab-forgejo-backup.service

[Install]
WantedBy=timers.target
